# service/discovery/run
#!/bin/sh
exec 2>&1
mkdir -p /usr/share/elasticsearch/config
exec ./discover.sh
